# software/init_scripts/infnoise.service
[Unit]
Description=Infinite Noise

[Service]
Type=simple
ExecStart=/usr/local/sbin/infnoise --dev-random

[Install]
WantedBy=multi-user.target

